# tweet-shots-api.service
[Unit]
Description=tweet-shots API Server
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/tweet-shots
ExecStart=/usr/bin/node api-server.mjs
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tweet-shots-api

# Environment variables
Environment=NODE_ENV=production
Environment=PORT=3000
# Set your admin key for managing API keys
Environment=ADMIN_KEY=your-secure-admin-key-here
# Optional: Set for URL response mode
# Environment=PUBLIC_URL=https://api.yourdomain.com/images
# Environment=OUTPUT_DIR=/var/lib/tweet-shots/output

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/tweet-shots

[Install]
WantedBy=multi-user.target
